# deploy/gradesentinel-backup.service
[Unit]
Description=GradeSentinel SQLite backup
Documentation=https://github.com/nzulkaynarov/GradeSentinel_bot
After=gradesentinel-bot.service

[Service]
Type=oneshot
User=gradesentinel
Group=gradesentinel
# .backup использует online-API SQLite — корректно копирует БД в режиме WAL
# без остановки бота. Хвост sed обрезает старые файлы (rotation 7 дней).
ExecStart=/bin/sh -c 'set -e; \
    ts=$(date -u +%%Y%%m%%d-%%H%%M%%S); \
    out=/var/backups/gradesentinel/sentinel-$ts.db; \
    /usr/bin/sqlite3 /var/lib/gradesentinel/sentinel.db ".backup $out"; \
    /usr/bin/gzip -9 $out; \
    find /var/backups/gradesentinel -name "sentinel-*.db.gz" -mtime +7 -delete; \
    echo "Backup OK: $out.gz ($(stat -c %%s $out.gz) bytes)"'

# Защита: оба пути требуют write — sqlite3 .backup может писать WAL/SHM
# на source для checkpoint'а, плюс ReadWritePaths для самого backup'а.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/gradesentinel /var/backups/gradesentinel
ProtectKernelTunables=true
ProtectKernelModules=true
RestrictSUIDSGID=true
