# deploy/gradesentinel-backup.timer
[Unit]
Description=Daily GradeSentinel SQLite backup
Documentation=https://github.com/nzulkaynarov/GradeSentinel_bot

[Timer]
# 03:30 Asia/Tashkent (UTC+5) → 22:30 UTC. Bot quiet hours (22:00–07:00 TST),
# минимум активности, низкий риск конкурренции с monitor'ом.
OnCalendar=*-*-* 22:30:00 UTC
Persistent=true
Unit=gradesentinel-backup.service

[Install]
WantedBy=timers.target
