# deploy/gradesentinel-webapp.service
[Unit]
Description=GradeSentinel WebApp (Flask)
Documentation=https://github.com/anthropics/gradesentinel
After=network-online.target gradesentinel-bot.service
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=gradesentinel
Group=gradesentinel
WorkingDirectory=/opt/gradesentinel
EnvironmentFile=/etc/gradesentinel/bot.env
Environment="DATABASE_PATH=/var/lib/gradesentinel/sentinel.db"
Environment="GOOGLE_CREDENTIALS_PATH=/etc/gradesentinel/credentials.json"
Environment="WEBAPP_PORT=8443"
Environment="WEBAPP_HOST=127.0.0.1"
Environment="PYTHONPATH=/opt/gradesentinel"
# gunicorn вместо Werkzeug dev server — multi-threaded, проде-grade.
# 2 worker x 4 threads = 8 одновременных запросов на 2 vCPU VPS.
# --timeout 30 — для медленных Anthropic вызовов в будущем (сейчас не используется в webapp).
# --access-logfile - выключаем (Caddy уже логирует), --error-logfile — в stderr (journalctl).
ExecStart=/opt/gradesentinel/venv/bin/gunicorn \
    --bind 127.0.0.1:8443 \
    --workers 2 \
    --threads 4 \
    --worker-class gthread \
    --timeout 30 \
    --keep-alive 5 \
    --access-logfile - \
    --error-logfile - \
    --log-level warning \
    webapp.app:app
Restart=always
RestartSec=10

# ── Hardening ────────────────────────────────────────────
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
# WebApp пишет в БД через тот же sqlite (init_db создаёт схему если нет).
ReadWritePaths=/var/lib/gradesentinel
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

# ── Лимиты ресурсов ──
MemoryMax=200M
MemoryHigh=180M
TasksMax=32

[Install]
WantedBy=multi-user.target
