# deploy/gradesentinel-offsite-backup.service
[Unit]
Description=GradeSentinel off-site backup sync (rclone → cloud)
Documentation=https://github.com/nzulkaynarov/GradeSentinel_bot
# После локального бэкапа (он пишет свежий .db.gz в 22:30 UTC) и поднятия сети.
After=gradesentinel-backup.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=gradesentinel
Group=gradesentinel
# rclone-кэш/логи в writable HOME (а не в root-only ~).
Environment=HOME=/var/lib/gradesentinel
ExecStart=/opt/gradesentinel/deploy/offsite-backup.sh

# Hardening: читаем только бэкапы + конфиг rclone, пишем только свой HOME/cache.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadOnlyPaths=/var/backups/gradesentinel /etc/gradesentinel
ReadWritePaths=/var/lib/gradesentinel
ProtectKernelTunables=true
ProtectKernelModules=true
RestrictSUIDSGID=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
