# deploy/gradesentinel-offsite-backup.timer
[Unit]
Description=Daily GradeSentinel off-site backup sync
Documentation=https://github.com/nzulkaynarov/GradeSentinel_bot

[Timer]
# 23:00 UTC = 04:00 Asia/Tashkent — через 30 мин после локального бэкапа
# (22:30 UTC), чтобы синкать уже свежий .db.gz. Тихие часы бота, низкая нагрузка.
OnCalendar=*-*-* 23:00:00 UTC
Persistent=true
Unit=gradesentinel-offsite-backup.service

[Install]
WantedBy=timers.target
