# deploy/gradesentinel-heartbeat.timer
[Unit]
Description=Run GradeSentinel heartbeat check every minute
Requires=gradesentinel-heartbeat.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=10s
Persistent=true

[Install]
WantedBy=timers.target
